# extras/arancino.service
[Unit]
Description=Arancino
After=network.target, redis-persistent.service, redis-volatile.service

[Service]
Type=simple
User=root
Group=root
ExecStart=<PATH TO ARANCINO MODULE>/start.py
Restart=on-failure
RestartSec=3
#on-abort

[Install]
WantedBy=multi-user.target
